# backend/crontab
0 9-17 * * 1-5 python /code/manage.py send_subscription_stock_emails_scheduled